# 99-razer-blackwidow-2014.rules
ACTION=="add", SUBSYSTEM=="input", ATTRS{idVendor}=="1532", ATTRS{idProduct}=="011b", RUN+="/usr/bin/razer-blackwidow-2014.py"
